# apps/sandbox-runtime/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/sandbox/supervisord.log
pidfile=/run/sandbox/supervisord.pid

[program:xvfb]
command=/bin/bash /opt/sandbox/scripts/start_x.sh
priority=10
autorestart=true
stdout_logfile=/var/log/sandbox/xvfb.log
stderr_logfile=/var/log/sandbox/xvfb.err.log

[program:chromium]
command=/bin/bash /opt/sandbox/scripts/start_browser.sh
priority=20
autorestart=true
stdout_logfile=/var/log/sandbox/chromium.log
stderr_logfile=/var/log/sandbox/chromium.err.log
environment=DISPLAY=":99"

[program:vnc]
command=/bin/bash /opt/sandbox/scripts/start_vnc.sh
priority=30
autorestart=true
stdout_logfile=/var/log/sandbox/vnc.log
stderr_logfile=/var/log/sandbox/vnc.err.log
environment=DISPLAY=":99"

